# systemd/lianli-daemon.service
[Unit]
Description=Lian Li Device Daemon
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/lianli-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
